# sratim.service
[Unit]
Description=Sratim Media Server
After=network.target

[Service]
Type=simple
User=sratim
WorkingDirectory=/var/lib/sratim
ExecStart=/usr/local/bin/sratim
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
